# svc/ptydaemon/ptydaemon.service.template
[Unit]
Description=Omni PTY daemon for %i
After=network.target

[Service]
Type=simple
User=%i
ExecStart=${PTYDAEMON_BIN}
Restart=on-failure
RestartSec=3s

# systemd creates /run/omni-%i/ owned by %i:root (mode 0700)
RuntimeDirectory=omni-%i
RuntimeDirectoryMode=0700

# systemd creates /var/lib/omni-%i/ for the persistent DB
StateDirectory=omni-%i
StateDirectoryMode=0700

Environment=OMNI_PTY_SOCKET=/run/omni-%i/omni-pty.sock
Environment=PTYDAEMON_DB=/var/lib/omni-%i/ptydaemon.db
# Include ~/.local/bin so user-local agent CLIs (claude, codex, gemini) are found.
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/%i/.local/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
